# digicast.service
[Unit]
Description=DigiCast
After=multi-user.target

[Service]
ExecStart=/home/digicast/run.sh

[Install]
WantedBy=multi-user.target